# rpm/after-install
#!/bin/sh

systemctl enable nginx
systemctl enable nexus
systemctl enable crl
